# s6-overlay/s6-rc.d/fail2ban/run
#!/command/execlineb -P
fail2ban-server -xf start
